# dnscrypt-systemd-units/dnscrypt-proxy@securedns-ipv6.socket
[Unit]
Description=dnscrypt-proxy listening socket

[Socket]
ListenStream=127.0.0.1:53539
ListenDatagram=127.0.0.1:53539

[Install]
WantedBy=sockets.target
